# supervisord-service.ini
[program:VAR_SERVICE_NAME]
command=VAR_NODE_PATH index.js
autostart=yes
startsecs=2
startretries=5
autorestart=yes
directory=VAR_DIRECTORY
environment=BOT_TOKEN="VAR_BOT_TOKEN"
